# systemd/llama-server.service
[Unit]
Description=Local Qwen 0.5B llama-server
After=network.target

[Service]
Type=simple
User=@USER@
Group=@USER@
ExecStart=/opt/llama.cpp/build/bin/llama-server \
  -m /opt/models/qwen2.5-0.5b-instruct-q4_k_m.gguf \
  --alias qwen2.5-0.5b-instruct \
  --host 127.0.0.1 \
  --port 8080 \
  -c 1024 \
  -t 4 \
  -tb 4 \
  -ngl 0 \
  -n 160
Restart=on-failure
RestartSec=5
Nice=5

[Install]
WantedBy=multi-user.target
